# src/networking/LeaderBoard/LeaderBoardApi.conf
#!upstart
description "LeaderBoardApi"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env NODE_ENV=production

exec node ~/networking/Api.js >> /var/log/LeaderBoardApi.log 2>&1